# packaging/debian/hosts-dns-forwarder.service
[Unit]
Description=hosts-dns-forwarder (serve /etc/hosts, otherwise forward)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hosts-dns-forwarder
Group=hosts-dns-forwarder

EnvironmentFile=-/etc/hosts-dns-forwarder/hosts-dns-forwarder.env
ExecStart=/usr/lib/hosts-dns-forwarder/hosts-dns-forwarder-run

AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=true

Restart=on-failure
RestartSec=2s

# Basic hardening
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictAddressFamilies=AF_INET AF_INET6
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
